# etc/systemd/system/alsa-restore.service.d/override.conf
[Unit]
Requires=dev-snd-controlC0.device
After=dev-snd-controlC0.device

[Service]
ExecStartPost=/usr/bin/amixer set Headphone -- -18dB
